# demo/crontab
0 0,12 * * * /home/braster/mONiToR/demo/restart.sh
